# deploy/configs/bdtree.service
[Unit]
Description=BDTree Family Tree Application
After=network.target

[Service]
Type=notify
User=bdtree
Group=bdtree
WorkingDirectory=/opt/bdtree
Environment="PATH=/opt/bdtree/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/bdtree/.env
ExecStart=/opt/bdtree/venv/bin/gunicorn -c /opt/bdtree/gunicorn.conf.py app:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
